# sound-eval/sound-eval.service
[Unit]
Description=Sound_evaluation_utility_for_exosensepi

[Service]
Type=simple
ExecStart=/usr/local/bin/sound-eval.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
